# .config/systemd/user/dunst.service
[Unit]
Description=dunst for notifications

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target
